# deployment/vigirrd.cron.in
# Remove old cached graphs
*/5 * * * *	@HTTPD_USER@	vigirrd-cleanup-cache
