# deployment/vigirrd.cron
# Remove old cached graphs
*/5 * * * *	apache	vigirrd-cleanup-cache
