# deployment/vigirrd.cron
# Remove old cached graphs
*/5 * * * *	root	vigirrd-cleanup-cache
